# build/etc/services.d/docker-stderr/run
#!/bin/sh

exec tail -f /docker-stderr >&2
