# build/etc/services.d/crond/run
#!/usr/bin/with-contenv sh

export PATH="$PATH:/var/www/html"
exec crond -L /dev/stdout -f
