# build/etc/services.d/crond/run
#!/usr/bin/with-contenv sh

exec crond -L /dev/stdout -f
