# build/etc/services.d/log/run
#!/bin/sh

exec tail -f /tmp/yii-stdout
